# cron/com.csnl.paper-rec.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.csnl.paper-rec.weekly</string>

  <!-- Run every Friday at 14:00 (Asia/Seoul; TZ env var forces it). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>5</integer>
    <key>Hour</key><integer>14</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/csnl/Documents/claude/csnl-paper-rec</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/csnl/Documents/claude/csnl-paper-rec/scripts/run_weekly_cron.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>Asia/Seoul</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/csnl/Documents/claude/csnl-paper-rec/state/cron.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/csnl/Documents/claude/csnl-paper-rec/state/cron.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
